# lab9/mephi.service
[Unit]
Description=Custom PHP Server
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=tecatech
ExecStart=/usr/bin/env php /home/tecatech/server.php

[Install]
WantedBy=multi-user.target