# lab9/tecatech.service
[Unit]
Description=Custom PHP Server
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=tecatexas
ExecStart=/usr/bin/env php /home/tecatexas/server.php

[Install]
WantedBy=multi-user.target